# budget-tracker.service
[Unit]
Description=Budget Tracker API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/budget-tracker
Environment=NODE_ENV=production
Environment=PORT=4001
Environment=DATA_FILE=/opt/budget-tracker/budget_data.json
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

